# deploy/supervisord.conf
[program:shorter]
command=bin/shorter -port %(ENV_SHORTER_PORT)
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=logs/shorter.log
exitcodes=0
